# crontab
0 0 * * * docker exec -it ap1 python worker.py 2
0 0 * * * docker exec -it ap1 python worker_365.py